# resources/rc.local
#!/bin/sh
su ooklaserver -c './etc/ooklaserver/OoklaServer --daemon'
exit 0